# crontab
# Backup databases daily at midnight
0 0 * * * /root/cron_task.sh /root/db_backup.sh

